# srcs/influxdb/supervisor.conf
[supervisord]
nodaemon=true
user=root

[program:influx]
command = sh -c "influxd -config /etc/influxdb.conf && kill -s SIGTERM $(cat supervisord.pid)"
